# kant-pastebin.service
[Unit]
Description=Kant Pastebin - UUCP + zkTLS
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/data1/kant/pastebin
ExecStart=/nix/store/s5g66rgab3hj1gzj3c8f438j8z8njw95-kant-pastebin-0.1.0/bin/kant-pastebin
Restart=always
RestartSec=10
Environment="BIND_ADDR=127.0.0.1:8090"
Environment="UUCP_SPOOL=/mnt/data1/spool/uucp/pastebin"
Environment="BASE_PATH=/pastebin"
Environment="BASE_URL=https://solana.solfunmeme.com"
Environment="NFT_DIR=/mnt/data1/time-2026/03-march/13/nft_enriched"
Environment="RUST_LOG=info"
Environment="PATH=/home/mdupont/.nix-profile/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
